# dots/config/systemd/user/battery-notif.service
[Unit]
Description=Battery notification

[Service]
ExecStart=%h/.local/bin/battery.sh
